# 10-floppy.rules
KERNEL=="sda", RUN+="/home/pi/floppy/floppy.sh"
